# wallpaper-switch.service
[Unit]
Description=NASA daily wallpaper

[Service]
Type=oneshot
ExecStart=/a/path/to/wallpaper-switch
