# wallpaper-switch.timer
[Unit]
Description=Nasa daily wallaper timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
